# systemd/rpi-video-viewer.service
[Unit]
Description=RGB LED Matrix video viewer
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/rpi-video-viewer
ExecStart=/bin/bash -lc '/usr/local/bin/video-viewer --led-chain=2 --led-cols=64 --led-rows=32 --led-gpio-mapping=adafruit-hat -T2 /home/pi/py/optimized-videos'
Restart=always
RestartSec=2s
TimeoutStopSec=5s
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
